# ansible/todo-app.socket
[Unit]
Description=My To-Do App's server

[Socket]
ListenStream=80
